# rc.local

ethtool -K eth1 tso off
ethtool -K eth1 gso off
ethtool -K eth1 gro off

